# nasbot.service
# ═══════════════════════════════════════════════════════════════════
#  NASBot Systemd Service — Auto-recovery dopo crash/reboot
#  Per TerraMaster F2-212
# ═══════════════════════════════════════════════════════════════════
#
#  INSTALLAZIONE:
#  1. Copia questo file in /etc/systemd/system/nasbot.service
#  2. Modifica Environment con i tuoi valori BOT_TOKEN e BOT_USER_ID
#  3. Esegui:
#     sudo systemctl daemon-reload
#     sudo systemctl enable nasbot
#     sudo systemctl start nasbot
#
#  COMANDI UTILI:
#     sudo systemctl status nasbot     # Stato
#     sudo systemctl restart nasbot    # Riavvia
#     sudo journalctl -u nasbot -f     # Log in tempo reale
#
# ═══════════════════════════════════════════════════════════════════

[Unit]
Description=NASBot - Telegram Bot per TerraMaster NAS
Documentation=https://github.com/tuouser/nasbot
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/Volume1/public

# ─── CONFIGURA QUESTI VALORI ───
Environment="BOT_TOKEN=IL_TUO_TOKEN_QUI"
Environment="BOT_USER_ID=IL_TUO_USER_ID_QUI"
# ───────────────────────────────

ExecStart=/Volume1/public/nasbot
ExecReload=/bin/kill -HUP $MAINPID

# Auto-restart in caso di crash
Restart=always
RestartSec=10

# Limiti e timeout
TimeoutStartSec=30
TimeoutStopSec=30
WatchdogSec=300

# Logging
StandardOutput=append:/Volume1/public/nasbot.log
StandardError=append:/Volume1/public/nasbot.log

# Sicurezza (opzionale, commenta se causa problemi)
# NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths=/Volume1/public

[Install]
WantedBy=multi-user.target
